# supervisor/supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/home/user/supervisord.pid
nodaemon=true
user=root
strip_ansi=false
childlogdir=/home/user/logs

[program:nginx]
command=/bin/bash -lc '/usr/local/openresty/bin/openresty -p /home/user/nginx -c /home/user/nginx/nginx.conf -g "daemon off;"'
directory=/home/user/nginx
autostart=true
autorestart=true
priority=5
stopsignal=TERM
stdout_logfile=/home/user/logs/nginx.supervisor.out
stderr_logfile=/home/user/logs/nginx.supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:syncd]
command=/bin/bash -lc 'if [ -n "${GITHUB_REPO:-}" ] && [ -n "${GITHUB_PAT:-}" ]; then exec /home/user/.venv/bin/python -m sync; else echo "[syncd] GITHUB_REPO/GITHUB_PAT not set; sync disabled"; exec sleep infinity; fi'
directory=/home/user
autostart=true
autorestart=true
startsecs=1
priority=10
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:prepare-runtime]
command=/bin/bash -lc '/home/user/scripts/prepare-runtime.sh'
directory=/home/user
autostart=true
autorestart=false
startsecs=0
priority=15
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:postgres]
command=/bin/bash -lc '/home/user/scripts/run-postgres.sh'
directory=/home/user
autostart=true
autorestart=true
startsecs=2
priority=20
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:redis]
command=/bin/bash -lc '/home/user/scripts/run-redis.sh'
directory=/home/user
user=1000
autostart=true
autorestart=true
startsecs=2
priority=21
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:sub2api-restore]
command=/bin/bash -lc '/home/user/scripts/restore-sub2api-backup.sh'
directory=/home/user
autostart=true
autorestart=false
startsecs=0
priority=25
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:sub2api]
command=/bin/bash -lc '/home/user/scripts/run-sub2api.sh'
directory=/home/user
user=1000
autostart=true
autorestart=true
startsecs=2
priority=30
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:sub2api-backup]
command=/bin/bash -lc '/home/user/scripts/backup-sub2api.sh'
directory=/home/user
user=1000
autostart=true
autorestart=true
startsecs=1
priority=35
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:pgadmin4]
command=/bin/bash -lc 'PYTHONDONTWRITEBYTECODE=1 exec /home/user/scripts/run-pgadmin.sh'
directory=/home/user
user=1000
autostart=true
autorestart=true
startsecs=2
priority=40
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:filebrowser]
command=/bin/bash -lc '/home/user/scripts/run-filebrowser.sh'
directory=/home/user
user=1000
autostart=true
autorestart=true
priority=45
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:gotty]
command=/bin/bash -lc '. /home/user/scripts/common-env.sh; wait_runtime_prepared; exec /home/user/gotty --address 127.0.0.1 --port 8081 --path /t --credential "${ADMIN_EMAIL}:${ADMIN_PASSWORD}" -w bash'
directory=/home/user
user=1000
autostart=true
autorestart=true
priority=50
stopsignal=TERM
stdout_logfile=/home/user/logs/gotty.supervisor.out
stderr_logfile=/home/user/logs/gotty.supervisor.err
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:cloudflared]
command=/bin/bash -lc '/home/user/scripts/run-cloudflared.sh'
directory=/home/user
autostart=true
autorestart=true
priority=60
stopsignal=TERM
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
